# job/launchd/com.blackjtsai.bjspg.post.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.blackjtsai.bjspg.post</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>/Users/blackjtsai/Documents/blackj-stock-point-guard/job/run_analysis.sh POST</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>21</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/blackjtsai/Documents/blackj-stock-point-guard/job/logs/launchd_post.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/blackjtsai/Documents/blackj-stock-point-guard/job/logs/launchd_post.err.log</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
